# conf.d/darkice.conf
[program:Darkice]
command=darkice
autostart=true
autorestart=true
priority=1
stderr_logfile=/var/log/supervisor/darkice.log
stdout_logfile=/var/log/supervisor/darkice.log
